# ioproxy.service
[Unit]
Description=IOProxy
After=network.target

[Service]
ExecStart=/usr/sbin/ioproxy

[Install]
WantedBy=multi-user.target